# .config/autostart/ksuperkey.desktop
[Desktop Entry]
Exec=ksuperkey
Icon=system-run
Path=
Terminal=false
Type=Application
